# app-utils/acme.sh/autobuild/overrides/usr/lib/systemd/user/acme-sh.service
[Unit]
Description=acme.sh cron job for daily certificate renewal

[Service]
ExecStart=/usr/bin/sh /usr/bin/acme.sh --cron --home %h/.acme.sh
