# app-database/postgresql-18/02-server/postgresqld.service.in
[Unit]
Description=PostgreSQL database server Version @MAJORVER@
After=network.target

[Service]
Type=forking
TimeoutSec=120
User=postgres
Group=postgres

Environment=PGROOT=/var/lib/postgres

SyslogIdentifier=postgres
PIDFile=/var/lib/postgres/data-@MAJORVER@/postmaster.pid
RuntimeDirectory=postgresql
RuntimeDirectoryMode=755

ExecStartPre=/usr/lib/postgresql-@MAJORVER@/bin/postgresql-check-db-dir ${PGROOT}/data-@MAJORVER@
ExecStart= /usr/lib/postgresql-@MAJORVER@/bin/pg_ctl -s -D ${PGROOT}/data-@MAJORVER@ start -w -t 120
ExecReload=/usr/lib/postgresql-@MAJORVER@/bin/pg_ctl -s -D ${PGROOT}/data-@MAJORVER@ reload
ExecStop=  /usr/lib/postgresql-@MAJORVER@/bin/pg_ctl -s -D ${PGROOT}/data-@MAJORVER@ stop -m fast

# Due to PostgreSQL's use of shared memory, OOM killer is often overzealous in
# killing Postgres, so adjust it downward
OOMScoreAdjust=-200

# Additional security-related features
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

